# test/Microsoft.AspNetCore.Server.Kestrel.FunctionalTests/SystemdActivation/activate-kestrel.socket
[Unit]
PartOf=activate-kestrel.service

[Socket]
ListenStream=8080
NoDelay=true

[Install]
WantedBy=multi-user.target
