# test/Microsoft.AspNetCore.Server.Kestrel.FunctionalTests/SystemdActivation/activate-kestrel.socket
[Unit]
Description=Kestrel Activation

[Socket]
ListenStream=8080
NoDelay=true

